# barorder.conf
[program:barorder]
command=/bin/bash /root/barorder/run.sh
autostart=true
autorestart=true
stopasgroup=true
stderr_logfile=/var/log/barorder.err.log
stdout_logfile=/var/log/barorder.out.log
